# django/cron/crontab
MAILTO=ninx@paperearth.net

* * * * * ~/gitrepos/paperearth.net/django/scripts/cronscript
