# tutorial/cron/crontab
MAILTO=ninx@paperearth.net

* * * * * ~/gitrepos/paperearth.net/django-tutorial/scripts/cronscript
